# tools/launchd/news-api.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>__LAUNCHD_PREFIX__.news-api</string>
    <key>ProgramArguments</key>
    <array>
        <string>__USER_HOME__/Development/nucleus/target/release/news-api</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__USER_HOME__/Development/nucleus</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>TZ</key>
        <string>__TZ__</string>
        <key>NUCLEUS_LOG</key>
        <string>info</string>
        <key>NUCLEUS_CLAUDE_BIN</key>
        <string>__USER_HOME__/.local/bin/claude</string>
        <key>PATH</key>
        <string>__USER_HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__USER_HOME__/Development/nucleus/memory/news-api.log</string>
    <key>StandardErrorPath</key>
    <string>__USER_HOME__/Development/nucleus/memory/news-api.log</string>
</dict>
</plist>
